# puppet/modules/custom/profiles/templates/mailcatcher/upstart.conf.erb
#!upstart

description "Mailcatcher"
start on runlevel [23]
stop on shutdown

# run with -f to be in foreground + bind on 0.0.0.0 to allow vagrant port forwarding
exec sudo -u vagrant bash -c "source /home/vagrant/.rvm/scripts/rvm && mailcatcher --http-ip 0.0.0.0 -f"

respawn