# puppet/modules/custom/profiles/files/php/php-fpm.conf
start on runlevel [2345]
stop on runlevel [016]

pre-start script
mkdir -p /run/php
chmod 0755 /run/php
chown www-data:www-data /run/php
end script

respawn
exec /usr/sbin/php-fpm --nodaemonize